# build/valhall.service
[Unit]
Description=aggregator
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/aggregator/bin
ExecStart=/bin/sh /opt/aggregator/bin/bootstrap.sh

User=svip
Group=svip
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target